# launchagent/com.pink.squid-pet.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pink.squid-pet</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT__/.venv/bin/python</string>
        <string>-m</string>
        <string>squid_pet</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>__PROJECT__/src</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/squid-pet.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/squid-pet.err.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
